# supervisor.conf
[supervisord]
user=root
nodaemon=true

[program:cron]
command = crond -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cron.err.log
stdout_logfile=/var/log/cron.log
